# tools/launchd/com.sovereign.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sovereign</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__/apps/sovereign</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/apps/sovereign/scripts/launcher/sovereign-daemon.sh</string>
  </array>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/Sovereign/sovereign.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/Sovereign/sovereign.log</string>
</dict>
</plist>
